# assignement_4_MP3_player/mp3_overlay/etc/profile.d/postbuild.sh
#!/bin/sh
hello.o
cd /myApplications
./sound_player.sh
